# root/etc/services.d/glusterfs/run
#!/usr/bin/with-contenv bash

exec 2>&1 \
    s6-setuidgid ${PUID:-911}:${PGID:-911} /usr/sbin/glusterd
