# root/etc/services.d/glusterfs/run
#!/usr/bin/with-contenv bash

exec 2>&1 \
    s6-setuidgid ${PUID:-1000}:${PGID:-1000} /usr/sbin/glusterfsd
